# alarm.cron
*/10 * * * * bash check_battery.sh
0 0 * * * echo false > data/low_battery_notification_sent
0 6 * * * kill -2 $(cat data/alarm.pid)
0 15 * * * bash proximity.sh
#0 7 * * * bash enable_wifi.sh false
#0 15 * * * bash enable_wifi.sh true