# alarm.cron
*/10 * * * * bash check_battery.sh
0 0 * * * echo false > data/low_battery_notification_sent
* 15 * * * bash proximity.sh